# examples/supervisor/conf.d/terminator.conf
[program:terminator]
command=/home/pi/terminator/bin/start
directory=/home/pi/terminator/bin
autostart=true
autorestart=true
user=pi
