# examples/supervisor/conf.d/web.conf
[program:web]
command=/home/pi/terminator/bin/start-web
directory=/home/pi/terminator/bin
autostart=false
autorestart=true
user=pi
